# pi-setup/alfredo-wake.service
[Unit]
Description=Alfredo Wake Word Listener
After=network.target alfredo-bridge.service
Wants=alfredo-bridge.service

[Service]
Type=simple
User=todd
WorkingDirectory=/home/todd/alfredo-kiosk
ExecStart=/usr/bin/python3 /home/todd/alfredo-kiosk/alfredo-wake.py
Restart=always
RestartSec=5
EnvironmentFile=/etc/alfredo/alfredo-kiosk.env
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/home/todd/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
