# pi-setup/alfredo-self-heal.service
[Unit]
Description=Alfredo Self-Heal — escalating recovery

[Service]
Type=oneshot
User=todd
ExecStart=/home/todd/alfredo-bridge/alfredo-self-heal.sh
