# pi-setup/alfredo-watchdog.timer
[Unit]
Description=Run Alfredo Watchdog every 2 minutes

[Timer]
OnBootSec=60
OnUnitActiveSec=120
AccuracySec=30

[Install]
WantedBy=timers.target
