# pi-setup/alfredo-self-heal.timer
[Unit]
Description=Run Alfredo Self-Heal every 5 minutes

[Timer]
OnBootSec=120
OnUnitActiveSec=300
AccuracySec=30

[Install]
WantedBy=timers.target
